# statsd-agent.conf
description "statsd-agent - Monitoring client"

start on runlevel [2345]
stop on [!2345] 

respawn

exec python /opt/statsd-agent/statsd-agent.py
